# artifacts/files/inittab
::restart:/bin/switch_root -c /dev/console ${_GOVLY} /init || ${ROOTFS_SHELL} -l

#>>> ::sysinit:${ROOTFS_SHELL} /initrc
#>>> ::respawn:-${ROOTFS_SHELL}
::respawn:${ROOTFS_SHELL} /initrc || ${ROOTFS_SHELL} -l
