# _artifacts/files/inittab
#>>> ::sysinit:${ROOTFS_SHELL} /initrc
#>>> ::respawn:-${ROOTFS_SHELL}
::respawn:${ROOTFS_SHELL} /initrc || ${ROOTFS_SHELL} -l

::restart:${ROOTFS_SROOT} -c /dev/console ${_GOVLY} /init || ${ROOTFS_SHELL} -l
